# scripts/com.hemnet-search.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent: runs the Hemnet update every other day.
  Install:
    cp scripts/com.hemnet-search.update.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.hemnet-search.update.plist
  A user LaunchAgent runs in your login session and does NOT need cron's
  Full Disk Access. Logs go to data/update.log (+ launchd.*.log).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hemnet-search.update</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/lihel5/playground/screening_hemnet/scripts/update.sh</string>
  </array>

  <!-- Fire daily at 12:30 (machine likely awake); update.sh gates to every
       other day. launchd runs a missed run once on the next wake. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>12</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/lihel5/playground/screening_hemnet/data/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/lihel5/playground/screening_hemnet/data/launchd.err.log</string>
</dict>
</plist>
